# airootfs/etc/profile.d/custom.sh
#!/bin/bash
systemctl enable NetworkManager.service
systemctl enable gdm.service
